# supervisord-haproxy.conf
[program:haproxy]
command=haproxy -db -f /etc/haproxy/haproxy.cfg
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supevisor/haproxy.log

