# supervisord-haproxy.conf
[program:haproxy]
command=haproxy -db -f /etc/haproxy/haproxy.cfg
numprocs=1
autostart=true
autorestart=true

