# supervisord-synapse.conf
[program:synapse]
command=/usr/local/bin/synapse -c /synapse.conf.json
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/synpase.log

